# conf/supervisor-nginx.conf
[program:nginx]
command=nginx -c /etc/nginx/nginx.conf
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

